# vmodem-x86.service
[Unit]
Description=VModem (x86)

[Service]
Type=simple
ExecStart=/usr/bin/vmodem_x86
