# vmodem-x86.service
[Unit]
Description=Virtual modem daemon for emulator

[Service]
Type=simple
ExecStart=/usr/bin/vmodem_x86
Environment=VGSM_VERBOSE=12
